# contrib/systemd/linearfs.service
[Unit]
Description=LinearFS - Linear issues as a FUSE filesystem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/linearfs/env
ExecStartPre=/usr/bin/mkdir -p ${LINEARFS_MOUNT}
ExecStart=%h/bin/linearfs mount -f ${LINEARFS_MOUNT}
ExecStop=/usr/bin/fusermount3 -u ${LINEARFS_MOUNT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
